# .rules/100-monitor-connection.rules
KERNEL=="card0" SUBSYSTEM=="drm", RUN+="/home/lasanjin/.plugs/monitor"