# .rules/101-usb-connected.rules
#ACTION=="add", ATTRS{idVendor}=="0781", ATTRS{idProduct}=="5583", RUN+="/home/lasanjin/.plugs/<script>"
