# .rules/101-usb-connection.rules
#ACTION=="add", ATTRS{idVendor}=="0781", ATTRS{idProduct}=="5583", RUN+="<PATH>"