# .rules/103-mouse-connected.rules
KERNEL=="hci0:3585", SUBSYSTEM=="bluetooth", ACTION=="add", RUN+="/home/lasanjin/.plugs/bluetooth m"