# .rules/104-keyboard-connected.rules
KERNEL=="hci0:256", SUBSYSTEM=="bluetooth", ACTION=="add", RUN+="/home/lasanjin/.plugs/bluetooth k"